# website.service
[Unit]
Description=website

[Service]
ExecStart=/home/pi/website/app.py
Restart=always

[Install]
WantedBy=default.target
